# mcp-server.service
[Unit]
Description=MCP Server
After=network.target

[Service]
User=your_username
WorkingDirectory=/path/to/AutoGen
Environment="PATH=/path/to/AutoGen/venv/bin"
EnvironmentFile=/path/to/AutoGen/.env
ExecStart=/path/to/AutoGen/venv/bin/gunicorn --bind 0.0.0.0:5000 "src.mcp_server.server:app"
Restart=always

[Install]
WantedBy=multi-user.target 